# deployment/khyber-shawls.service
[Unit]
Description=Khyber Shawls Next.js Application
Documentation=https://github.com/webspires-eng/khyber-shawls
After=network.target mysql.service

[Service]
Type=simple
User=deploy
WorkingDirectory=/home/deploy/khyber-shawls
Environment=NODE_ENV=production
EnvironmentFile=/home/deploy/khyber-shawls/.env
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=khyber-shawls

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/deploy/khyber-shawls/public/uploads

[Install]
WantedBy=multi-user.target
